# supervisor/citserver.conf
[program:citserver]
command=/usr/local/bin/startcitserver
autostart=true
autorestart=true
priority=1
startsecs=5
startretries=1000000
