# supervisor/webcithttp.conf
[program:webcithttp]
command=/usr/local/bin/startwebcit http
autostart=true
autorestart=true
startsecs=5
startretries=1000000
