# supervisor/webcithttps.conf
[program:webcithttps]
command=/usr/local/bin/startwebcit https
autostart=true
autorestart=true
startsecs=5
startretries=1000000
